# supervisor-local.conf
[program:local]
command=/etc/rc.local
autorestart=false
startsecs=0
